# systemd/callscoot-agent.service
[Unit]
Description=CallScoot AI call agent
After=pipewire.service pipewire-pulse.service wireplumber.service callscoot-daemon.service
Wants=pipewire.service pipewire-pulse.service wireplumber.service

[Service]
Type=simple
EnvironmentFile=-%h/.config/callscoot/elevenagents.env
ExecStart=%h/.local/lib/callscoot/.venv/bin/python %h/.local/lib/callscoot/callscoot_agent.py run --mode elevenagents
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
